# ops/launchagents/com.hongstr.weekly_backfill.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.hongstr.weekly_backfill</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-lc</string>
        <string>cd __REPO_ROOT__ &amp;&amp; source scripts/load_env.sh &amp;&amp; bash scripts/backfill_1m_from_2020.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>0</integer>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>PY</key>
        <string>__REPO_ROOT__/.venv/bin/python</string>
    </dict>
    <key>StandardOutPath</key>
    <string>__REPO_ROOT__/logs/launchd_weekly_backfill.out.log</string>
    <key>StandardErrorPath</key>
    <string>__REPO_ROOT__/logs/launchd_weekly_backfill.err.log</string>
</dict>
</plist>
