# ops/launchagents/com.hongstr.tg_cp_watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.hongstr.tg_cp_watchdog</string>

    <key>ProgramArguments</key>
    <array>
      <string>/bin/bash</string>
      <string>-lc</string>
      <string>cd "/Users/hong/Projects/HONGSTR" &amp;&amp; ( [ -x ./.venv/bin/python ] &amp;&amp; ./.venv/bin/python scripts/tg_cp_healthcheck.py || python3 scripts/tg_cp_healthcheck.py ) || launchctl kickstart -k gui/$(id -u)/com.hongstr.tg_cp</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
      <key>TG_CP_MAX_AGE_SEC</key>
      <string>7200</string>
    </dict>

    <key>StartInterval</key>
    <integer>120</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/hong/Library/Logs/hongstr/tg_cp_watchdog.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/hong/Library/Logs/hongstr/tg_cp_watchdog.err.log</string>
  </dict>
</plist>
