# ops/launchagents/com.hongstr.dashboard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.hongstr.dashboard</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>WorkingDirectory</key>
    <string>__REPO_ROOT__</string>

    <key>ProgramArguments</key>
    <array>
      <string>/bin/bash</string>
      <string>-lc</string>
      <string>cd "__REPO_ROOT__" &amp;&amp; (source scripts/load_env.sh || true) &amp;&amp; exec .venv/bin/streamlit run scripts/dashboard.py --server.address 127.0.0.1 --server.port 8501 --server.headless true</string>
    </array>

    <key>StandardOutPath</key>
    <string>__REPO_ROOT__/logs/launchd_dashboard.out.log</string>
    <key>StandardErrorPath</key>
    <string>__REPO_ROOT__/logs/launchd_dashboard.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
      <key>PY</key>
      <string>__REPO_ROOT__/.venv/bin/python</string>
      <key>HONGSTR_HOST</key>
      <string>127.0.0.1</string>
      <key>HONGSTR_PORT</key>
      <string>8501</string>
    </dict>
  </dict>
</plist>
