# config/com.hongstr.research_poller.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.hongstr.research_poller</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/hong/Projects/HONGSTR/scripts/poll_research_loop.sh</string>
    </array>

    <!-- Every 10 minutes -->
    <key>StartInterval</key>
    <integer>600</integer>

    <!-- Run once on load to catch any queued triggers immediately -->
    <key>RunAtLoad</key>
    <true/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>PYTHONPATH</key>
        <string>/Users/hong/Projects/HONGSTR</string>
    </dict>

    <key>WorkingDirectory</key>
    <string>/Users/hong/Projects/HONGSTR</string>

    <key>StandardOutPath</key>
    <string>/Users/hong/Projects/HONGSTR/logs/launchd_research_poller.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/hong/Projects/HONGSTR/logs/launchd_research_poller.err.log</string>

    <key>KeepAlive</key>
    <false/>
</dict>
</plist>
